# launchd/com.stockanalyze.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- 每日盤後(14:30)更新母體股價 + 重算指標。載入見 launchd/README。 -->
<plist version="1.0">
<dict>
  <key>Label</key>          <string>com.stockanalyze.daily</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/kaochihping/Stock_analyze/scripts/run_daily.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>   <integer>14</integer>
    <key>Minute</key> <integer>30</integer>
  </dict>
  <key>RunAtLoad</key>          <false/>
  <key>StandardOutPath</key>    <string>/Users/kaochihping/Stock_analyze/logs/launchd.daily.out</string>
  <key>StandardErrorPath</key>  <string>/Users/kaochihping/Stock_analyze/logs/launchd.daily.err</string>
  <key>WorkingDirectory</key>   <string>/Users/kaochihping/Stock_analyze</string>
</dict>
</plist>
